# rc.local
#rc.local
#loop
fi
	python /home/pi/Documents/TripWire/main.py &
exit 0
